# etc/supervisord/supervisord.conf
; Sample supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Notes:
;  - Shell expansion ("~" or "$HOME") is not supported.  Environment
;    variables can be expanded using this syntax: "%(ENV_HOME)s".
;  - Comments must have a leading space: "a=b ;comment" not "a=b;comment".

[unix_http_server]
file=/home/dataman/supervisord/supervisor.sock
chmod=0700                                       ; socket file mode (default 0700)

[supervisord]
logfile=/home/dataman/supervisord/log/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/home/dataman/supervisord/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
childlogdir=/home/dataman/supervisord/log         ; ('AUTO' child log dir, default $TEMP)
user=dataman
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
nodaemon=true               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///home/dataman/supervisord/supervisor.sock   ; use a unix:// URL  for a unix socket

[include]
files = /home/dataman/supervisord/programs/*.conf

[program:n6log_std]
command=/usr/bin/tail -f /home/dataman/logs/log_n6_all
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
autorestart=true
startsecs=3
